# com.leegonzales.research-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leegonzales.research-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/leegonzales/Projects/miniconda/bin/python3</string>
        <string>-m</string>
        <string>research_agent.cli.main</string>
        <string>run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/leegonzales/Projects/leegonzales/Intelbot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/leegonzales/Projects/miniconda/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/leegonzales/Projects/leegonzales/Intelbot</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leegonzales/.research-agent/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/leegonzales/.research-agent/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
